# infra/atlas-host-executor.service
[Unit]
Description=Atlas Host-Executor Bridge
Documentation=https://github.com/GainMarketing123/nanoclaw
After=nanoclaw.service
Wants=nanoclaw.service

[Service]
Type=simple
User=atlas
Group=atlas
WorkingDirectory=/home/atlas/nanoclaw
ExecStart=/usr/bin/python3 /home/atlas/nanoclaw/host/host-executor.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=atlas-host-executor

# Environment
Environment=HOME=/home/atlas
Environment=TZ=America/New_York
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/home/atlas/.local/bin
# Explicit Atlas state + Nanoclaw repo paths so future User= migrations don't
# silently re-resolve via the new user's $HOME. host-executor.py reads these
# at module load with `Path.home() / ...` defaults if unset, so deployments
# that haven't picked up the new unit (or that haven't run daemon-reload
# after the unit lands) continue to behave as before.
Environment=ATLAS_DIR=/home/atlas/.atlas
Environment=NANOCLAW_DIR=/home/atlas/nanoclaw

# Minimal hardening — claude -p needs broad filesystem access for hooks
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
