# infra/atlas-mission-control.service
[Unit]
Description=Atlas Mission Control Dashboard
Documentation=https://github.com/GainMarketing123/nanoclaw
After=nanoclaw.service

[Service]
Type=simple
User=atlas
Group=atlas
WorkingDirectory=/home/atlas/nanoclaw
ExecStart=/usr/bin/node /home/atlas/nanoclaw/mission-control/server.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=atlas-mission-control

# Environment
Environment=HOME=/home/atlas
Environment=TZ=America/New_York
Environment=MC_PORT=8080

[Install]
WantedBy=multi-user.target
