# infra/atlas-host-executor.service
[Unit]
Description=Atlas Host-Executor Bridge
Documentation=https://github.com/GainMarketing123/nanoclaw
After=nanoclaw.service
Wants=nanoclaw.service

[Service]
Type=simple
User=atlas
Group=atlas
WorkingDirectory=/home/atlas/nanoclaw
ExecStart=/usr/bin/python3 /home/atlas/nanoclaw/host/host-executor.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=atlas-host-executor

# Environment
Environment=HOME=/home/atlas
Environment=TZ=America/New_York
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/home/atlas/.local/bin

# Minimal hardening — claude -p needs broad filesystem access for hooks
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
